# etc/systemd/system/checkupdate.service
[Unit]
Description=run checkupdate every hour

[Service]
ExecStart=bash ~/scripts/updaterun

[Install]
WantedBy=checkupdate.target
