# etc/systemd/system/checkupdate.timer
[Unit]
Description=Run checkupdate every hour

[Timer]
OnCalendar=0/1:00:00
Persistent=true

[Install]
WantedBy=checkupdate.target